# scripts/codex-runner-dashboard.service.example
[Unit]
Description=Nexus Codex Runner Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nexus-codex-runner
EnvironmentFile=/home/ubuntu/nexus-codex-runner/dashboard.env
ExecStart=/usr/bin/node /home/ubuntu/nexus-codex-runner/codex-runner-dashboard.mjs
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
