# scripts/codex-runner@.service.example
[Unit]
Description=Nexus Codex Runner Instance %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nexus-codex-runner
EnvironmentFile=/home/ubuntu/nexus-codex-runner/.env
# Set the runner id in ExecStart so it wins over a shared CODEX_RUNNER_ID in .env.
ExecStart=/usr/bin/env CODEX_RUNNER_ID=%i /usr/bin/node /home/ubuntu/nexus-codex-runner/codex-runner.mjs
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
