# scripts/codex-runner@.service.example
[Unit]
Description=Nexus Codex Runner Instance %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nexus-codex-runner
EnvironmentFile=/home/ubuntu/nexus-codex-runner/.env
Environment=CODEX_RUNNER_ID=%i
ExecStart=/usr/bin/node /home/ubuntu/nexus-codex-runner/codex-runner.mjs
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
